# deploy/localagent.container
[Unit]
Description=localagent gateway
After=network-online.target

[Container]
Image=localhost/localagent:latest
ContainerName=localagent-gateway
Exec=gateway
PodmanArgs=--runtime=krun
UserNS=keep-id
Network=pasta
PublishPort=18790:18790
PublishPort=18791:18791
Volume=%h/.localagent:/home/localagent/.localagent:Z
AutoUpdate=local

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
